# cookbooks/rails-app/templates/ssh-forward-ec2.service.erb
[Unit]
Description=SshForwardEc2

[Service]
Type=simple
ExecStart=/home/pi/ssh_forward
User=pi
TimeoutSec=30
RestartSec=15s

[Install]
WantedBy=multi-user.target
